# ops/systemd/hr-hunter-restore-drill.service
[Unit]
Description=HR Hunter backup restore drill
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/hr-hunter/current
EnvironmentFile=/srv/hr-hunter/shared/env/hr-hunter.env
Environment=PYTHONPATH=/srv/hr-hunter/current/src
ExecStart=/bin/bash -lc 'set -euo pipefail; /srv/hr-hunter/.venv/bin/python -m hr_hunter.cli runtime-restore-drill --workspace-root /srv/hr-hunter --temp-root /srv/hr-hunter/shared/monitoring/restore-drills'

[Install]
WantedBy=multi-user.target
